# lxc/mascot.service
[Unit]
Description=KI mit Bedacht – KI-Maskottchen (Flask/Gunicorn)
After=network.target
Wants=network-online.target

[Service]
User=ki-mit-bedacht
Group=ki-mit-bedacht
WorkingDirectory=/opt/ki-mit-bedacht

# .env wird automatisch von python-dotenv geladen (WorkingDirectory)
ExecStart=/opt/ki-mit-bedacht/venv/bin/gunicorn \
    --worker-class=gthread \
    --workers=1 \
    --threads=4 \
    --bind=0.0.0.0:5000 \
    --timeout=120 \
    --access-logfile=- \
    --error-logfile=- \
    app:app

Restart=always
RestartSec=5

# Logs → journald (abrufbar mit: journalctl -u ki-mit-bedacht -f)
StandardOutput=journal
StandardError=journal
SyslogIdentifier=ki-mit-bedacht

# Hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ReadWritePaths=/opt/ki-mit-bedacht

[Install]
WantedBy=multi-user.target
